# infra/systemd/doit-reminders.service
[Unit]
Description=Doit Reminder Dispatcher (Production)
After=doit.service
Requires=doit.service

[Service]
Type=oneshot
User=salomao
Group=salomao
EnvironmentFile=/srv/doit/prod/doit-config/web.env
ExecStart=/bin/sh -c '/usr/bin/curl -fsS --max-time 30 -X POST -H "Authorization: Bearer $CRON_SECRET" http://127.0.0.1:8110/api/notifications/reminders'
StandardOutput=journal
StandardError=journal
