# infra/systemd/doit-dev.service
[Unit]
Description=Doit Web Application (Development)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=salomao
Group=salomao
WorkingDirectory=/srv/doit/dev/app
EnvironmentFile=/srv/doit/dev/doit-config/web.env
Environment=NODE_ENV=production
Environment=PORT=8111
Environment=HOSTNAME=127.0.0.1
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/usr/bin/corepack pnpm --filter @doit/web exec next start -H 127.0.0.1 -p 8111
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
